# Downloads/Spartan_Labs/website/spartan-agents.service
[Unit]
Description=Spartan Research Station - Autonomous Agent Swarm
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=spartan
WorkingDirectory=/mnt/c/Users/Quantum/Downloads/Spartan_Labs/website
ExecStart=/usr/bin/python3 /mnt/c/Users/Quantum/Downloads/Spartan_Labs/website/agent_orchestrator.py
Restart=always
RestartSec=10
StandardOutput=append:/mnt/c/Users/Quantum/Downloads/Spartan_Labs/website/agent_orchestrator.log
StandardError=append:/mnt/c/Users/Quantum/Downloads/Spartan_Labs/website/agent_orchestrator.log

# Environment
Environment="PYTHONUNBUFFERED=1"

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
